# deploy/systemd/sandboxd-controller.service
[Unit]
Description=Sandbox control plane
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
User=sandbox
Group=sandbox
EnvironmentFile=/etc/sandbox/sandbox.env
ExecStart=/usr/local/bin/sandboxd --role controller --config /etc/sandbox/sandbox.toml
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/sandbox
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
